# s6-rc.d/cron/run
#!/command/execlineb -P

exec crond -f -d 7
